# automounts/srv-projects.automount
[Unit]
Description=SMB Projects Automount
ConditionPathIsEmpty=/srv/projects

[Automount]
Where=/srv/projects
TimeoutIdleSec=0

[Install]
WantedBy=multi-user.target
